# cookbooks/openruko/templates/default/upstart-openruko-logplex.conf.erb
start on starting openruko
stop on stopping openruko
respawn

exec su - vagrant -c 'cd /home/vagrant/openruko/logplex; \
logplex/bin/logplex >> /var/log/openruko/logplex.log 2>&1'