# cookbooks/openruko/templates/default/upstart-openruko-gitmouth.conf.erb
start on starting openruko
stop on stopping openruko
respawn

exec su - vagrant -c 'cd /home/vagrant/openruko/gitmouth; ./bin/gitmouth >> /var/log/openruko/gitmouth.log 2>&1'
